# systemd/idle-tracker.service
[Unit]
Description=Idle Tracker

[Service]
ExecStart=idle_tracker.py
Restart=on-failure

[Install]
WantedBy=gnome-session.target
